# backup-service/crontab
# Ejecutar backup completo diariamente a las 2:00 AM
0 2 * * * /app/scripts/backup-daily.sh >> /var/log/backup.log 2>&1

# Ejecutar backup de BD cada 6 horas
0 */6 * * * /app/scripts/backup-db-only.sh >> /var/log/backup-db.log 2>&1
